# services/cloudflared/run
#!/data/data/com.termux/files/usr/bin/sh
exec cloudflared tunnel --config /data/data/com.termux/files/home/.cloudflared/config.yml run 2>&1
