# services/jellyfin/run
#!/data/data/com.termux/files/usr/bin/sh
exec /data/data/com.termux/files/usr/bin/jellyfin 2>&1
